# systemd/music-bot.service.j2
[Service]
Environment=TELEGRAM_TOKEN={{ 8137581059:AAFOSgGg1RE-Kq0coZceSYncj_L89DMdYK0 }}
ExecStart=/opt/app/.venv/bin/python -u bot.py
