# ansible/music-bot.service.j2
[Unit]
Description=Music Telegram Bot
After=network.target music-api.service
Requires=music-api.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/app
Environment=TELEGRAM_TOKEN={{ telegram_token }}
ExecStart=/opt/app/.venv/bin/python -u bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
